# docker/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv bash
# shellcheck shell=bash

set -e

echo "❯ Starting nginx ..."
exec nginx
