# docker/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv bash
set -e

echo "â¯ Starting nginx ..."
exec s6-setuidgid npmuser nginx
