# docker/rootfs/etc/s6-overlay/s6-rc.d/prepare/up
# shellcheck shell=bash
/etc/s6-overlay/s6-rc.d/prepare/script.sh
