# docker/rootfs/etc/s6-overlay/s6-rc.d/frontend/run
#!/command/with-contenv bash
# shellcheck shell=bash

set -e

# This service is DEVELOPMENT only.

if [ "$DEVELOPMENT" = 'true' ]; then
	. /usr/bin/common.sh
	cd /app/frontend || exit 1
	HOME=$NPMHOME
	export HOME
	mkdir -p /app/frontend/dist
	chown -R "$PUID:$PGID" /app/frontend/dist

	log_info 'Starting frontend ...'
	s6-setuidgid "$PUID:$PGID" yarn install
	exec s6-setuidgid "$PUID:$PGID" yarn watch
else
	exit 0
fi
