# docker/rootfs/etc/services.d/manager/finish
#!/usr/bin/with-contenv bash

s6-svscanctl -t /var/run/s6/services
