# rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bash

mkdir -p /tmp/nginx \
  /data/{nginx,logs,access} \
  /data/nginx/stream \
  /var/lib/nginx/cache/{public,private}

chown root /tmp/nginx
exec nginx

