# rootfs/etc/services.d/manager/run
#!/usr/bin/with-contenv bash

cd /srv/manager
node --abort_on_uncaught_exception --max_old_space_size=250 /srv/manager/src/backend/index.js
